# jobs/reaper/templates/config/reaper.cron
*/59 * * * * /var/vcap/jobs/reaper/bin/reaper >>/var/vcap/sys/log/reaper/reaper.log 2>&1
